# ibusdev.target
[Unit]
Description=I/K-bus network devices
RefuseManualStart=yes
